# wc2026-workers-bot.service
[Unit]  
Description=Telegram WC2026 Workers Prediction Bot  
After=network-online.target  
Wants=network-online.target  
  
[Service]  
Type=simple  
WorkingDirectory=/root/wc2026/workers/Telegram_WC2026_Prediction_Bot  
EnvironmentFile=/root/wc2026/workers/Telegram_WC2026_Prediction_Bot/.env  
ExecStart=/root/wc2026/workers/Telegram_WC2026_Prediction_Bot/.venv/bin/python /root/wc2026/workers/Telegram_WC2026_Prediction_Bot/main.py  

Restart=always  
RestartSec=10  
  
[Install]  
WantedBy=multi-user.target
